# config/.archived/cronjobs/bibliography.cron
*/5 * * * * $HOME/.dotfiles/scripts/notes/make-bibliography.sh