# scheduler/cronjobs/snts.cron
0 0 * * 1 $HOME/.dotfiles/scripts/archive-note-logs.sh
